# src/fty-kpi-power-uptime.service.in
[Unit]
Description=42ity kpi-power-uptime service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-kpi-power-uptime

[Install]
WantedBy=bios.target
