# src/fty-kpi-power-uptime.service.in
[Unit]
Description=42ity kpi-power-uptime service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-kpi-power-uptime 

[Install]
WantedBy=bios.target
